# config/php7-clean.timer
[Unit]
Description=Default PHP7 clean service session

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
